# bluetooth.rules
ACTION=="add", KERNEL=="bluetooth", RUN+="/bin/sh -c '[ -x /etc/init.d/bluetooth ] && /etc/init.d/bluetooth restart >/dev/null'"
